# roles/sync-configuration/templates/sync-configuration.service
[Unit]
Description=Run Ansible playbook which controls the universe.

[Service]
ExecStart=/usr/bin/make {{ make_target }}
WorkingDirectory={{ workdir | expanduser }}

[Install]
WantedBy=multi-user.target
