# roles/sync-dotfiles/templates/sync-dotfiles.timer
[Unit]
Description=Update git repository with dotfiles every day.

[Timer]
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
